# contrib/mold-server.service
[Unit]
Description=mold inference server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mold serve --port 7680 --bind 0.0.0.0
Restart=on-failure
RestartSec=5

# Environment — override via EnvironmentFile or drop-in
Environment=MOLD_LOG=info
Environment=MOLD_PORT=7680
# For secrets (HF_TOKEN, MOLD_API_KEY), use EnvironmentFile instead of inline Environment:
# EnvironmentFile=-/etc/mold/mold.env
# Then set HF_TOKEN=hf_... and MOLD_API_KEY=... in /etc/mold/mold.env (chmod 600)
#
# Optional authentication and rate limiting:
# Environment=MOLD_API_KEY=your-secret-key
# Environment=MOLD_RATE_LIMIT=10/min
# Environment=MOLD_RATE_LIMIT_BURST=20

# NixOS: CUDA driver access
# Environment=LD_LIBRARY_PATH=/run/opengl-driver/lib

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.mold %h/.config/mold %h/.cache/huggingface
PrivateTmp=true

[Install]
WantedBy=default.target
